# tic_scripts/bbb_disconnect.cron
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root /opt/greenlight/tic_scripts/bbb_signature_disconnect.sh
